# etc/supervisord.d/haproxy.ini
[program:haproxy]
command=/sbin/haproxy -db -f /etc/haproxy/haproxy.cfg
redirect_stderr=true
stdout_logfile = /var/log/supervisor/%(program_name)s.log

