# etc/supervisord.d/haproxy.ini
[program:haproxy]
command=/root/bin/start_haproxy.sh
redirect_stderr=true
stdout_logfile = /var/log/supervisor/%(program_name)s.log

